# deploy/batty.service
[Unit]
Description=Batty web UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/github/pi-face
Environment=BATTY_HOST=127.0.0.1
Environment=BATTY_PORT=3147
ExecStart=/usr/bin/env pnpm start -- /root/github
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
